# conf/cron/cron.d/skos-history-ysa-yso-update
# Update the skos-history database and YSO at 04:22 on the first day of each month
22 04 1 * * root /usr/local/bin/update-skos-history && /usr/local/bin/update-ysa-allars-to-yso
